# backup/backup.service
[Unit]
Description=Backup to S3

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup
